# services/sb-pd.service
[Unit]
Description=Sigmund blow Start PD patch

[Service]
ExecStart=/home/pi/src/sigmundBlow/scripts/sb-pd


[Install]
WantedBy=multi-user.target
