# services/sg-pd.service
[Unit]
Description=PSIS StartPD

[Service]
ExecStart=/home/pi/src/posteDecoute/scripts/psis-pd


[Install]
WantedBy=multi-user.target
